# deploy/ether-server.service
# Референс systemd-юнита для прода (etherapp.ru) — копия того, что стоит на
# сервере в /etc/systemd/system/ether-server.service. scripts/deploy.sh только
# рестартит уже установленный сервис. После правок: sudo systemctl daemon-reload.
#
# Раскладка раздельная: бинарник read-only в /opt/ether, конфиг в /etc/ether,
# изменяемое состояние (SQLite-база) — в /var/lib/ether. ProtectSystem=strict
# делает всю ФС read-only, а ReadWritePaths возвращает запись ТОЛЬКО в
# /var/lib/ether — то есть база отделена от бинарника и конфига намеренно.
[Unit]
Description=Ether server (geographic IRC backend)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ether
Group=ether
WorkingDirectory=/var/lib/ether
ExecStart=/opt/ether/ether-server -config /etc/ether/config.prod.json
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
ReadWritePaths=/var/lib/ether

[Install]
WantedBy=multi-user.target
